# hubo-ach_server/init/docpad-server.conf
description "The HTTP server for Hubo-in-the-Browser Live"

start on startup
stop on shutdown

script 
    cd /home/wash/drchubo.js
    su -c "docpad run" wash
end script
