# hubo-ach_server/init/socketio-server.conf
description "The Socket.IO data server for Hubo-in-the-Browser Live"

start on startup
stop on shutdown

script 
    cd /home/wash/drchubo.js/hubo-ach_server
    su -c "./run.sh" wash
end script
